# deploy/game-video-tool@.service
[Unit]
Description=Game Video Tool FastAPI Service on port %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/game-video-runtime/slots/%i/current
EnvironmentFile=/home/deploy/game-video-tool/.env
Environment=USER_DATA_DIR=/home/deploy/game-video-data
Environment=AUTH_ENABLED=true
Environment=UI_DIST_DIR=/home/deploy/game-video-runtime/slots/%i/current/react-ui/dist
ExecStart=/home/deploy/game-video-tool/.venv/bin/python server/main.py --port %i --host 127.0.0.1
Restart=always
RestartSec=3
TimeoutStopSec=15
StandardOutput=append:/home/deploy/game-video-runtime/logs/app-%i.log
StandardError=append:/home/deploy/game-video-runtime/logs/app-%i.log

[Install]
WantedBy=multi-user.target
